# community_relay/deploy/token-monitor-community.service
[Unit]
Description=Token Monitor Community Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=true
StateDirectory=token-monitor-community
StateDirectoryMode=0750
UMask=0027
EnvironmentFile=/etc/token-monitor-community.env
ExecStart=/usr/local/bin/token-monitor-community-relay
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/token-monitor-community
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
